# devel/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/home/arjan/Desktop/ros_noetic_base_2204/catkin_ws/src/ros/tools/rosbash/roszsh"
